# config/supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=info
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
priority=200